# roles/bh-rdiff-backup/templates/systemd/backup.service
[Unit]
Description=Backup hosts

[Service]
Type=oneshot
ExecStart=/home/backup/backup.sh
TimeoutSec=infinity

