# roles/bh-rdiff-backup/templates/systemd/backup.timer
[Unit]
Description=Run backups

[Timer]
# OnBootSec=15min
# OnUnitActiveSec=1d
OnCalendar=*-*-* 4:00:00

[Install]
WantedBy=multi-user.target
