# roles/bh-atomic/templates/base/borgbackup/borgbackup.service

[Unit]
Description=Backup hosts

[Service]
Type=oneshot
ExecStart=/usr/bin/sh -c /root/backup.sh
TimeoutSec=infinity


