# roles/bh-laptop/templates/laptop/rc.d/rc.local
#!/bin/sh
echo OFF > /sys/kernel/debug/vgaswitcheroo/switch

exit 0

