# roles-archive/bh-rdiff-backup/templates/systemd/backup.service
[Unit]
Description=Backup hosts

[Service]
Type=oneshot
ExecStart=/usr/bin/sh -c /home/backup/backup.sh
TimeoutSec=infinity

