# roles/bh-atomic/templates/zfs/zfs-auto-snapshot-hourly.service
[Unit]
Description=ZFS hourly snapshot service

[Service]
ExecStart=/root/zfs-auto-snapshot -d -k 72 tank
